# chip/etc/vdn/root.crontab
*/5 * * * * /usr/local/bin/fccid
*/10 * * * * /etc/vdn/network/./set_routing
* * * * * /etc/vdn/watchdog/./watchdog
*/15 * * * * /et/./.check_vdn
